# OS_Migration/services/mybot.service
[Unit]
Description=My ESP32 Bot Service
After=network.target
# Restart if network comes back after being down
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pi4_drive/Git_projects/RASPI4-MAIN

# Activate venv and run the script
ExecStart=/bin/bash -c 'source /home/pi/myenv/bin/activate && exec python /home/pi/pi4_drive/Git_projects/RASPI4-MAIN/main.py'

# Environment variables
Environment="PYTHONUNBUFFERED=1"
Environment="PATH=/home/pi/myenv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pi/.Xauthority"

# Logging: Python logging writes to stdout; telepot network timeout noise goes to stderr.
# Capture structured Python logs to the output file; suppress the raw telepot stderr flood.
StandardOutput=append:/home/pi/pi4_drive/Git_projects/RASPI4-MAIN/logs/Output_mybot_service.log
StandardError=append:/home/pi/pi4_drive/Git_projects/RASPI4-MAIN/logs/Error_mybot_service.log

# Restart policy — restart on any exit (crash, OOM, signal, etc.)
Restart=always
RestartSec=10

# If it crashes 5 times within 5 minutes → reboot the Pi
StartLimitBurst=5
StartLimitIntervalSec=300
OnFailure=systemd-reboot.service

[Install]
WantedBy=multi-user.target
